# deployment/sixpack-web.conf
[program:sixpack-web]
command=gunicorn -w 4 -b:8080 --worker-class=gevent sixpack.web:start
directory=/home/sixpack
user=sixpack
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0